# Linux-Repo - VK5-setup/mqtt-api.service
[Unit]
Description=MQTT Chat API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mqtt-chat
ExecStart=/home/ubuntu/mqtt-chat/api.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target