# Linux-Repo - VK5-setup/mqtt-logger.service
[Unit]
Description=MQTT to MySQL Logger
After=network.target mysql.service docker.service
Wants=mysql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mqtt-chat
ExecStart=/home/ubuntu/mqtt-chat/mqtt_logger.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
